# server-install/systemd/olcrtc-server.service
[Unit]
Description=olcRTC server (WebRTC tunnel through whitelisted services)
Documentation=https://github.com/openlibrecommunity/olcrtc
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
EnvironmentFile=/etc/olcrtc/env
User=olcrtc
Group=olcrtc
StateDirectory=olcrtc
StateDirectoryMode=0750
WorkingDirectory=/var/lib/olcrtc
ExecStart=/usr/local/bin/olcrtc \
    -mode srv \
    -provider ${OLCRTC_PROVIDER} \
    -id ${OLCRTC_ROOM_ID} \
    -key ${OLCRTC_KEY} \
    -dns ${OLCRTC_DNS}

# Hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallArchitectures=native
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=300
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
